# packaging/rpm/zmailer.cron
#!/bin/sh
##
## ZMailer default crontab:
##

# Resubmit deferred messages
28 */2 * * *	root	. /etc/zmailer/zmailer.conf ; $MAILBIN/zmailer resubmit 2>/dev/null 1>&2

# Cleanout public and postman directories
7 */2 * * *	root	. /etc/zmailer/zmailer.conf ; $MAILBIN/zmailer cleanup 2>/dev/null 1>&2
